# deployment/rabbithole/rabbithole.service
[Unit]
Description=Rabbit Hole Inc - Next.js App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/app/rabbithole/apps/rabbithole
ExecStart=/usr/bin/node /app/rabbithole/node_modules/next/dist/bin/next start -p 3000 -H 0.0.0.0
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
EnvironmentFile=-/etc/rabbithole/env

[Install]
WantedBy=multi-user.target
